# forms/detect_new.json
{
  "name": "New Detection Rule Nomination",
  "description": "Submit a Request for a New Detection Rule",
  "subject": "New Detection Rule Nomination",
  "elements": [
    {
      "label": "Detection System",
      "type": "input",
      "description": "Enter System that should  run the detection (eg. Firewall, NIDS, AV, etc.).",
      "required": true
    },
    {
      "label": "Describe the Rule",
      "type": "textarea",
      "placeholder": "Describe the circumstances when the rule should trigger",
      "description": "Describe when the Rule should Trigger.",
      "required": true,
      "rows": 5
    },
    {
      "label": "Urgency",
      "type": "select",
      "values": [
        "CRITICAL",
        "High",
        "Normal",
        "Routine"
      ],
      "default": "Routine",
      "description": "Select an Urgency Level for Implementation of this New Detection Rule.",
      "required": true
    },
    {
      "label": "Associated MITRE ATT&CK Tactic or Technique",
      "type": "select",
      "values": [
        " - ",
        "Reconnaissance (TA0043)",
        "Resource Development (TA0042)",
        "Initial Access (TA0001)",
        "Execution (TA0002)",
        "Persistence (TA0003)",
        "Privilege Escalation (TA0004)",
        "Defense Evasion (TA0005)",
        "Credential Access (TA0006)",
        "Discovery (TA0007)",
        "Lateral Movement (TA0008)",
        "Collection (TA0009)",
        "Command and Control (TA0011)",
        "Exfiltration (TA0010)",
        "Impact (TA0040)",
        "Access Token Manipulation (T1134)",
        "Account Discovery (T1087)",
        "Application Deployment Software (T1017)",
        "Application Layer Protocol (T1071)",
        "Automated Exfiltration (T1020)",
        "Boot or Logon Autostart Execution (T1547)",
        "Boot or Logon Initialization Scripts (T1037)",
        "Brute Force (T1110)",
        "Bypass User Account Control (T1548.002)",
        "Clipboard Data (T1115)",
        "Clipboard Modification (T1115)",
        "Command and Scripting Interpreter (T1059)",
        "Command-Line Interface (T1059.003)",
        "Command-Line Interface: PowerShell (T1059.001)",
        "Command-Line Interface: Windows Command Shell (T1059.003)",
        "Create Account (T1136)",
        "Create or Modify System Process (T1543)",
        "Create or Modify System Process: Launch Agent (T1543.001)",
        "Credentials in Files (T1552.001)",
        "Credentials in Registry (T1552.002)",
        "Custom Command and Control Protocol (T1094)",
        "Custom Cryptographic Protocol (T1024)",
        "Data Compressed (T1002)",
        "Data Destruction (T1485)",
        "Data Encrypted (T1022)",
        "Data Encrypted for Impact (T1486)",
        "Data from Cloud Storage Object (T1530)",
        "Data from Local System (T1005)",
        "Data from Network Shared Drive (T1039)",
        "Data Obfuscation (T1001)",
        "Data Staged (T1074)",
        "Data Transfer Size Limits (T1030)",
        "DLL Side-Loading (T1073)",
        "Disabling Security Tools (T1562)",
        "Email Collection (T1114)",
        "Exfiltration Over Alternative Protocol (T1048)",
        "Exfiltration Over Command and Control Channel (T1041)",
        "Execution through API (T1106)",
        "File and Directory Discovery (T1083)",
        "File and Directory Permissions Modification (T1222)",
        "File Deletion (T1070.004)",
        "Hijack Execution Flow (T1574)",
        "Hijack Execution Flow: DLL Search Order Hijacking (T1574.001)",
        "Hijack Execution Flow: DLL Side-Loading (T1073)",
        "Hijack Execution Flow: Shared Webroot (T1071.005)",
        "Inhibit System Recovery (T1490)",
        "Indicator Removal on Host (T1070)",
        "Indicator Removal on Host: Clear Linux or Mac System Logs (T1070.002)",
        "Indicator Removal on Host: Clear Windows Event Logs (T1070.001)",
        "Indicator Removal on Host: File Deletion (T1070.004)",
        "Input Capture (T1056)",
        "Input Capture: Credential API Hooking (T1056.004)",
        "Input Capture: Keyboard Input Capture (T1056.001)",
        "Ingress Tool Transfer (T1105)",
        "Logon Scripts (T1037.001)",
        "Masquerading (T1036)",
        "Modify Cloud Compute Infrastructure (T1578)",
        "Modify Registry (T1112)",
        "Network Service Scanning (T1046)",
        "Network Service Scanning: SMB/NetBios Scanning (T1046.001)",
        "Network Share Discovery (T1135)",
        "Network Sniffing (T1040)",
        "New Service (T1050)",
        "Obfuscated Files or Information (T1027)",
        "Phishing (T1566)",
        "Process Discovery (T1057)",
        "Process Injection (T1055)",
        "Process Injection: Dynamic-link Library Injection (T1055.001)",
        "Process Termination (T1489)",
        "Query Registry (T1012)",
        "Registry Run Keys / Startup Folder (T1547.001)",
        "Remote File Copy (T1105)",
        "Remote Services (T1021)",
        "Remote System Discovery (T1018)",
        "Remote System Discovery: Windows Admin Shares (T1018.001)",
        "Replication Through Removable Media (T1091)",
        "Scheduled Task (T1053)",
        "Scheduled Task/Job (T1053)",
        "Scheduled Task/Job: Launch Agent (T1053.001)",
        "Screen Capture (T1113)",
        "Security Software Discovery (T1518.001)",
        "Service Execution (T1035)",
        "Shared Modules (T1129)",
        "Signed Binary Proxy Execution (T1218)",
        "Signed Script Proxy Execution (T1216)",
        "Software Deployment Tools (T1072)",
        "Spearphishing Attachment (T1566.001)",
        "Sudo (T1169)",
        "Sudo and Sudo Caching (T1169)",
        "System Information Discovery (T1082)",
        "System Owner/User Discovery (T1033)",
        "System Service Discovery (T1007)",
        "System Time Discovery (T1124)",
        "Timestomp (T1070.006)",
        "User Execution (T1204)",
        "User Execution: Malicious File (T1204.002)",
        "User Execution: Malicious Link (T1204.001)",
        "Valid Accounts (T1078)",
        "Valid Accounts: Cloud Accounts (T1078.003)",
        "Valid Accounts: Domain Accounts (T1078.002)",
        "Valid Accounts: Local Accounts (T1078.001)",
        "Web Service (T1102)"
      ],
      "description": "Select the Applicable MITRE ATT&CK Technique (optional).",
      "default": " - "
    },
    {
      "label": "Suggested Duration of Detection Rule",
      "type": "input",
      "description": "Describe when the Rule should be Reviewed.",
      "default": "90 days"
    },
    {
      "label": "Justification",
      "type": "textarea",
      "placeholder": "Add comments supporting your request",
      "description": "Enter a Justification or Reason for the New Rule.",
      "required": true,
      "rows": 10
    }
  ]
}
